# etc/systemd/system/km-ru-mail.conf
description "SSL wrapper (km.ru)"

start on local-filesystems and net-device-up IFACE!=lo
stop on runlevel [016]

console log

exec /usr/bin/stunnel4 /etc/stunnel/km-ru-smtps.init
