# etc/systemd/system/yandex-mail.conf
description "SSL wrapper (Yandex)"

start on local-filesystems and net-device-up IFACE!=lo
stop on runlevel [016]

console log

exec /usr/bin/stunnel4 /etc/stunnel/yandex-smtps.init
